# var/lib/dpkg/info/bluetooth-touch.list
/.
/var
/var/lib
/var/lib/lxc
/var/lib/lxc/android
/var/lib/lxc/android/pre-start.d
/var/lib/lxc/android/pre-start.d/40-brcm_patchram_plus-disable
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bluetooth-touch.service
/usr
/usr/bin
/usr/bin/brcm_patchram_plus
/usr/share
/usr/share/bluetooth-touch
/usr/share/bluetooth-touch/flo
/usr/share/bluetooth-touch/maguro
/usr/share/bluetooth-touch/grouper
/usr/share/bluetooth-touch/mako
/usr/share/doc
/usr/share/doc/bluetooth-touch
/usr/share/doc/bluetooth-touch/changelog.Debian.gz
/usr/share/doc/bluetooth-touch/copyright
/etc
/etc/init
/etc/init/bluetooth-touch-android.conf
/etc/init/bluetooth-touch.conf
/etc/init/bluetooth-touch-maguro.conf
/etc/init/bluetooth-touch-grouper.conf
/etc/init/bluetooth-touch-flo.conf
/etc/init/bluetooth-touch-mako.conf
