# var/lib/dpkg/info/click.list
/.
/etc
/etc/click
/etc/click/databases
/etc/click/databases/20_custom.conf
/etc/click/databases/10_core.conf
/etc/click/databases/99_default.conf
/etc/init
/etc/init/click-system-hooks.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/click-system-hooks.service
/usr
/usr/bin
/usr/bin/click
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/click
/usr/lib/arm-linux-gnueabihf/click/libclickpreload.so
/usr/lib/systemd
/usr/lib/systemd/user
/usr/lib/systemd/user/click-user-hooks.service
/usr/share
/usr/share/mime
/usr/share/mime/packages
/usr/share/mime/packages/click.xml
/usr/share/lintian/overrides/click
/usr/share/upstart
/usr/share/upstart/sessions
/usr/share/upstart/sessions/click-user-hooks.conf
/usr/share/man/man1/click.1.gz
/usr/share/click
/usr/share/click/hooks
/usr/share/click/hooks/click-desktop.hook
/usr/share/doc
/usr/share/doc/click
/usr/share/doc/click/changelog.gz
/usr/share/doc/click/copyright
