# var/lib/dpkg/info/hybris-usb.list
/.
/etc
/etc/init
/etc/init/usb-tethering.conf
/etc/hybris-usb
/etc/hybris-usb/dhcpd.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/usb-tethering.service
/usr
/usr/bin
/usr/bin/usb-tethering
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/hybris-usb.conf
/usr/share
/usr/share/doc
/usr/share/doc/hybris-usb
/usr/share/doc/hybris-usb/changelog.gz
/usr/share/doc/hybris-usb/copyright
