# var/lib/dpkg/info/lxc-android-config.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/wpa_supplicant_touch.conf
/etc/system-image
/etc/system-image/writable-paths
/etc/logrotate.d
/etc/logrotate.d/touch-syslog
/etc/cron.hourly
/etc/cron.hourly/logrotate
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/50-firmware.rules
/etc/udev/rules.d/90-alsa-ucm.rules
/etc/udev/rules.d/90-alsa-restore.rules
/etc/udev/rules.d/60-persistent-v4l.rules
/etc/init
/etc/init/isc-dhcp-server.override
/etc/init/avahi-daemon.override
/etc/init/ssh.override
/etc/init/force-adb.conf
/etc/init/tty4.override
/etc/init/apport-reconfig.conf
/etc/init/cups-browsed.override
/etc/init/lxc-android-boot.conf
/etc/init/android-set-boot-completion-property.conf
/etc/init/alsa-restore.override
/etc/init/device-hacks.conf
/etc/init/cups.override
/etc/init/upstart-local-bridge.conf
/etc/init/apport-config.conf
/etc/init/apport.override
/etc/init/tty3.override
/etc/init/lxc-android-config.conf
/etc/init/boot-hooks
/etc/init/boot-hooks/wipe-qml-cache.conf
/etc/init/boot-hooks/extrausers.conf
/etc/init/boot-hooks/ensure-logdir-owner.conf
/etc/init/isc-dhcp-server6.override
/etc/init/rild.conf
/etc/init/override-vibrator-permission.conf
/etc/init/udev.override
/etc/init/tty5.override
/etc/init/tty1.override
/etc/init/override-leds-permission.conf
/etc/init/plymouth-upstart-bridge.override
/etc/init/ofono.override
/etc/init/alsa-store.override
/etc/init/no-cpu-hotplug.conf
/etc/init/ureadahead-other.override
/etc/init/force-mtp.conf
/etc/init/ureadahead.override
/etc/init/lxc-net.override
/etc/init/whoopsie.override
/etc/init/adbd-emergency-shell.conf
/etc/init/android-usb-state.conf
/etc/init/update-machine-info.conf
/etc/init/tethering.conf
/etc/init/boot-hooks-emit.conf
/etc/init/ureadahead-touch.conf
/etc/init/ssh-keygen.conf
/etc/init/setvtrgb.override
/etc/init/tty6.override
/etc/init/enable-cpu-hotplugging.conf
/etc/init/tty2.override
/etc/init/ubuntu-location-service.override
/etc/init/ssh-property-watcher.conf
/etc/init/zram-touch.conf
/etc/profile.d
/etc/profile.d/disable-qtwebengine-gpu.sh
/etc/rsyslog.d
/etc/rsyslog.d/45-touch.conf
/etc/NetworkManager
/etc/NetworkManager/dispatcher.d
/etc/NetworkManager/dispatcher.d/02default_route_workaround
/etc/NetworkManager/dispatcher.d/03mmsproxy
/etc/systemd
/etc/systemd/network
/etc/systemd/network/00-rmnet-kernel.link
/etc/ubuntu-touch-session.d
/etc/ubuntu-touch-session.d/tf700t.conf
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/70-android.rules
/lib/udev/rules.d/99-z-video.rules
/lib/udev/rules.d/85-android.rules
/lib/udev/rules.d/65-android.rules
/lib/udev/rules.d/99-android.rules
/lib/systemd
/lib/systemd/system
/lib/systemd/system/alsa-state.service.d
/lib/systemd/system/alsa-state.service.d/disabled.conf
/lib/systemd/system/force-mtp.service
/lib/systemd/system/alsa-restore.service.d
/lib/systemd/system/alsa-restore.service.d/disabled.conf
/lib/systemd/system/ofono.service.d
/lib/systemd/system/ofono.service.d/lxc-android-config.conf
/lib/systemd/system/upstart-local-bridge.service
/lib/systemd/system/rild.service
/lib/systemd/system/lxc-android-config.service
/var
/var/lib
/var/lib/lxc-android-config
/var/lib/lxc-android-config/diversions
/var/lib/lxc
/var/lib/lxc/android
/var/lib/lxc/android/rootfs
/var/lib/lxc/android/pre-start.d
/var/lib/lxc/android/pre-start.d/10-no-adbd
/var/lib/lxc/android/pre-start.d/25-process-overrides
/var/lib/lxc/android/pre-start.d/30-no-surface-flinger
/var/lib/lxc/android/pre-start.d/20-no-forced-adb
/var/lib/lxc/android/pre-start.d/50-micshm
/var/lib/lxc/android/pre-start.d/15-no-uchroot
/var/lib/lxc/android/pre-start.d/20-no-default-prop-adb
/var/lib/lxc/android/pre-start.sh
/var/lib/lxc/android/overrides
/var/lib/lxc/android/config
/usr
/usr/bin
/usr/bin/ssh-keygen.sh
/usr/bin/tethering
/usr/lib
/usr/lib/lxc-android-config
/usr/lib/lxc-android-config/70-flo.rules
/usr/lib/lxc-android-config/70-mako.rules
/usr/lib/lxc-android-config/70-nicki.rules
/usr/lib/lxc-android-config/70-hammerhead.rules
/usr/lib/lxc-android-config/70-GT-I9100.rules
/usr/lib/lxc-android-config/update-fstab
/usr/lib/lxc-android-config/70-u9200.rules
/usr/lib/lxc-android-config/70-generic.rules
/usr/lib/lxc-android-config/copy-udev-unflipped
/usr/lib/lxc-android-config/70-pme.rules
/usr/lib/lxc-android-config/70-tf700t.rules
/usr/lib/lxc-android-config/70-pollux_windy.rules
/usr/lib/lxc-android-config/70-pollux.rules
/usr/lib/lxc-android-config/70-deb.rules
/usr/lib/lxc-android-config/70-rolex.rules
/usr/lib/lxc-android-config/70-GT-N8013.rules
/usr/lib/lxc-android-config/70-GT-N8000.rules
/usr/lib/lxc-android-config/lxc-android-ready
/usr/lib/lxc-android-config/70-oneplus2.rules
/usr/lib/lxc-android-config/70-trunk.rules
/usr/lib/lxc-android-config/70-land.rules
/usr/lib/lxc-android-config/70-generic_x86.rules
/usr/lib/lxc-android-config/70-maguro.rules
/usr/lib/lxc-android-config/70-oneplus3.rules
/usr/lib/lxc-android-config/70-manta.rules
/usr/lib/lxc-android-config/70-grouper.rules
/usr/share
/usr/share/apparmor
/usr/share/apparmor/hardware
/usr/share/apparmor/hardware/audio.d
/usr/share/apparmor/hardware/audio.d/apparmor-easyprof-ubuntu_tf700t
/usr/share/apparmor/hardware/audio.d/apparmor-easyprof-ubuntu_android
/usr/share/apparmor/hardware/graphics.d
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_tf700t
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_k95v2
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_GT-N8000
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_krillin
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_GT-N8013
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_pollux_windy
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_turbo
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_pollux
/usr/share/apparmor/hardware/graphics.d/apparmor-easyprof-ubuntu_android
/usr/share/apparmor/hardware/video.d
/usr/share/apparmor/hardware/video.d/apparmor-easyprof-ubuntu_krillin
/usr/share/apparmor/hardware/video.d/apparmor-easyprof-ubuntu_android
/usr/share/upstart
/usr/share/upstart/sessions
/usr/share/upstart/sessions/android-usb-state.conf
/usr/share/powerd
/usr/share/powerd/device_configs
/usr/share/powerd/device_configs/config-tf700t.xml
/usr/share/doc
/usr/share/doc/lxc-android-config
/usr/share/doc/lxc-android-config/changelog.gz
/usr/share/doc/lxc-android-config/copyright
