# var/lib/dpkg/info/ofono.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/ofono.conf
/etc/ofono
/etc/ofono/ril_subscription.conf
/etc/ofono/phonesim.conf
/etc/init
/etc/init/ofono.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ofono.service
/var
/var/lib
/var/lib/ofono
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/ofono
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ofonod.8.gz
/usr/share/doc
/usr/share/doc/ofono
/usr/share/doc/ofono/AUTHORS.gz
/usr/share/doc/ofono/changelog.gz
/usr/share/doc/ofono/copyright
/usr/sbin
/usr/sbin/ofonod
