# var/lib/dpkg/info/packagekit.list
/.
/etc
/etc/PackageKit
/etc/PackageKit/PackageKit.conf
/etc/PackageKit/events
/etc/PackageKit/events/pre-transaction.d
/etc/PackageKit/events/pre-transaction.d/README
/etc/PackageKit/events/post-transaction.d
/etc/PackageKit/events/post-transaction.d/README
/etc/PackageKit/Vendor.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.PackageKit.conf
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/20packagekit
/usr
/usr/share
/usr/share/PackageKit
/usr/share/PackageKit/pk-upgrade-distro.sh
/usr/share/PackageKit/helpers
/usr/share/PackageKit/helpers/test_spawn
/usr/share/PackageKit/helpers/test_spawn/search-name.sh
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.PackageKit.Transaction.xml
/usr/share/dbus-1/interfaces/org.freedesktop.PackageKit.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.PackageKit.service
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.packagekit.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.packagekit.rules
/usr/share/doc
/usr/share/doc/packagekit
/usr/share/doc/packagekit/copyright
/usr/share/doc/packagekit/NEWS.gz
/usr/share/doc/packagekit/README.Debian
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/packagekit-plugins
/usr/lib/arm-linux-gnueabihf/packagekit-plugins/libpk_plugin-scan-desktop-files.so
/usr/lib/arm-linux-gnueabihf/packagekit-plugins/libpk_plugin_scripts.so
/usr/lib/arm-linux-gnueabihf/packagekit-plugins/libpk_plugin-update-package-cache.so
/usr/lib/arm-linux-gnueabihf/packagekit-plugins/libpk_plugin-systemd-updates.so
/usr/lib/arm-linux-gnueabihf/packagekit-backend
/usr/lib/arm-linux-gnueabihf/packagekit-backend/libpk_backend_test_spawn.so
/usr/lib/arm-linux-gnueabihf/packagekit-backend/libpk_backend_dummy.so
/usr/lib/arm-linux-gnueabihf/packagekit-backend/libpk_backend_test_thread.so
/usr/lib/arm-linux-gnueabihf/packagekit-backend/libpk_backend_test_succeed.so
/usr/lib/arm-linux-gnueabihf/packagekit-backend/libpk_backend_test_nop.so
/usr/lib/arm-linux-gnueabihf/packagekit-backend/libpk_backend_test_fail.so
/usr/lib/packagekit
/usr/lib/packagekit/pk-clear-offline-update
/usr/lib/packagekit/packagekitd
/usr/lib/packagekit/pk-offline-update
/usr/lib/packagekit/pk-trigger-offline-update
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/packagekit.service
/lib/systemd/system/packagekit-offline-update.service
/var
/var/cache
/var/cache/PackageKit
/var/cache/PackageKit/downloads
/var/lib
/var/lib/PackageKit
/usr/share/doc/packagekit/changelog.Debian.gz
