# var/lib/dpkg/info/ubuntu-touch-session.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/ofono-phablet.conf
/etc/pulse
/etc/pulse/touch.pa
/etc/pulse/touch-stream-restore.table
/etc/pulse/touch-client.conf
/etc/ubuntu-platform-api
/etc/ubuntu-platform-api/application.conf
/etc/init
/etc/init/lightdm.override
/etc/init/plymouth.override
/etc/init/plymouth-shutdown.override
/etc/init/plymouth-log.override
/etc/init/keep-display-on-during-boot.conf
/etc/init/factory-wipe.conf
/etc/init/no-android.conf
/etc/init/apparmor.override
/etc/profile.d
/etc/profile.d/ubuntu-touch-file-selector.sh
/etc/profile.d/qml_qv4_jit_cache.sh
/etc/profile.d/add_custom_to_xdg_data.sh
/etc/profile.d/qt_networkmanager_bearer.sh
/etc/profile.d/dbus-source.sh
/etc/profile.d/upstart-phablet.sh
/etc/ubuntu-touch-session.d
/etc/ubuntu-touch-session.d/GT-N8000.conf
/etc/ubuntu-touch-session.d/flo.conf
/etc/ubuntu-touch-session.d/cdma_targa.conf
/etc/ubuntu-touch-session.d/manta.conf
/etc/ubuntu-touch-session.d/maguro.conf
/etc/ubuntu-touch-session.d/tf300t.conf
/etc/ubuntu-touch-session.d/umts_spyder.conf
/etc/ubuntu-touch-session.d/u9200.conf
/etc/ubuntu-touch-session.d/SGH-I777.conf
/etc/ubuntu-touch-session.d/espressowifi.conf
/etc/ubuntu-touch-session.d/deb.conf
/etc/ubuntu-touch-session.d/generic_x86.conf
/etc/ubuntu-touch-session.d/GT-N8013.conf
/etc/ubuntu-touch-session.d/cdma_maserati.conf
/etc/ubuntu-touch-session.d/cdma_solana.conf
/etc/ubuntu-touch-session.d/dlx.conf
/etc/ubuntu-touch-session.d/espressorf.conf
/etc/ubuntu-touch-session.d/grouper.conf
/etc/ubuntu-touch-session.d/hammerhead.conf
/etc/ubuntu-touch-session.d/GT-I9100.conf
/etc/ubuntu-touch-session.d/cdma_spyder.conf
/etc/ubuntu-touch-session.d/GT-I9000B.conf
/etc/ubuntu-touch-session.d/mako.conf
/etc/ubuntu-touch-session.d/generic.conf
/etc/ubuntu-touch-session.d/android.conf
/etc/ubuntu-touch-session.d/GT-N7000.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/display-manager.service.d
/lib/systemd/system/display-manager.service.d/ubuntu-touch-session.conf
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/com.ubuntu.touch.NetworkManager.pkla
/usr
/usr/bin
/usr/bin/create-online-accounts
/usr/share
/usr/share/upstart
/usr/share/upstart/sessions
/usr/share/upstart/sessions/bluez-mpris-proxy.conf
/usr/share/upstart/sessions/obexd.conf
/usr/share/upstart/sessions/xdg-dirs.conf
/usr/share/upstart/sessions/indicator-sound.override
/usr/share/upstart/sessions/ofono-setup.conf
/usr/share/upstart/sessions/pulseaudio-trust-stored.conf
/usr/share/upstart/sessions/cameraservice-trust-stored.conf
/usr/share/upstart/sessions/pulseaudio.conf
/usr/share/doc
/usr/share/doc/ubuntu-touch-session
/usr/share/doc/ubuntu-touch-session/changelog.gz
/usr/share/doc/ubuntu-touch-session/copyright
/usr/share/pam-configs
/usr/share/pam-configs/touch-extrausers
/usr/sbin
/usr/sbin/ubuntu-touch-lightdm-session
